# crontab
# Just a dummy
*   *   *   *   *   echo "Please provide a crontab - example job: sync example-config.conf"